# crontab
*/15 * * * * ubuntu bash /home/ubuntu/apps/changing-lives/repository/cron.sh fab $DEPLOYMENT_TARGET deploy_aggregates
